# crontab
0,30 * * * * speedtest --csv >> ~/dotfiles/speedtest.csv
